# test-crontab.cron
# Run script every day at 2 AM
0 2 * * * /media/sf_jenkins-project/gs-spring-boot/cron-backup.sh
